# files/etc/profile.d/04-isengard-vim.sh
# This alias is to override the default vimrc file with the one from isengard-vim
alias vim="vim -u /isengard-vim/vimrc"