# files/etc/profile.d/05-isengard-neovim.sh
# This alias is to override the default vimrc file with nvchad from isengard-neovim
alias nvim="nvim -u /isengard-neovim/nvim"
