# files/etc/profile.d/99-isengard-env.sh
export EDITOR=vim

eval "$(zoxide init bash)"
