# files/etc/profile.d/99-isengard-env.sh
export EDITOR=nvim

eval "$(zoxide init bash)"
eval "$(gh copilot alias -- bash)"
