# files/etc/profile.d/99-isengard-env.sh
export EDITOR=emacs

eval "$(zoxide init bash)"

PATH=$PATH:/XDG_DIRS/config/emacs/bin
