# files/etc/profile.d/isengard-atuin.sh
bind -x '"\C-p": __atuin_history --shell-up-key-binding'
